# mt5/supervisord.conf
; =============================================================================
; Supervisor config for MT5 Container
; Manages: MT5 Terminal (Wine), ZMQ bridge (Wine Python)
; Xvfb is started by startup.sh before supervisor.
; =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid

[supervisorctl]
serverurl=unix:///var/run/supervisor/supervisor.sock

[unix_http_server]
file=/var/run/supervisor/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ---------------------------------------------------------------------------
; MT5 Terminal -- runs inside Wine
; ---------------------------------------------------------------------------
[program:mt5_terminal]
command=/bin/bash /app/run_mt5.sh
autostart=true
autorestart=true
startsecs=30
startretries=5
priority=20

; ---------------------------------------------------------------------------
; Wine Python rpyc server -- exposes MT5 API via rpyc on port 18812
; ---------------------------------------------------------------------------
[program:mt5_rpyc_server]
command=/bin/bash /app/run_rpyc.sh
autostart=true
autorestart=true
startsecs=10
startretries=5
priority=30

; ---------------------------------------------------------------------------
; ZeroMQ Bridge -- runs in Wine Python, translates ZMQ to MT5 directly
; ---------------------------------------------------------------------------
[program:zmq_bridge]
command=wine python /app/zmq_bridge.py
autostart=true
autorestart=true
startsecs=5
startretries=5
priority=40
